# systemd/genelec-sleep.service
[Unit]
Description=Genelec Sleep
DefaultDependencies=no
Before=shutdown.target
Conflicts=reboot.target

[Service]
Type=oneshot
ExecStart=GenelecCli sleep 
TimeoutSec=3

[Install]
WantedBy=shutdown.target
