# systemd/genelec-wake.service
[Unit]
Description=Genelec Wake

[Service]
ExecStart=genelec-cli wake

[Install]
WantedBy=multi-user.target
